# .config/autostart/gromit-mpx.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/gromit-mpx
