# .config/autostart/qterminal.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/qterminal -d
